# backup-daemon.service
[Unit]
Description=Backup Daemon

[Service]
ExecStart=/home/leo/Desktop/daemon/build/Backup
Restart=always

[Install]
WantedBy=default.target
